# backend/docker/com.ascoder.tls-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ascoder.tls-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Developer/CommandLineTools/usr/bin/python3</string>
        <string>/Users/weigaolei/CodeSpace/WorkSpace/qys-skill/Ascoder/git-tls-proxy.py</string>
        <string>8443</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/ascoder-tls-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ascoder-tls-proxy.log</string>
</dict>
</plist>
